# deployment/gunicorn_supervisor.conf
[program:APP_NAME]
command=PROJECT_DIR/bin/gunicorn_start
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=PROJECT_DIR/logs/gunicorn-error.log
